# services/runFullReconciliation.timer
[Unit]
Description=Run Full Reconciliation Pipeline periodically
After=network.target neo4j.service

[Timer]
# Run every 6 hours
OnBootSec=5min
OnUnitActiveSec=6h
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
